# timer.timer
[Unit]
Requires=task_to_backup.service

[Timer]
Unit=myBackuper.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
