# timer.timer
[Unit]
Requires=Mybackuper.service

[Timer]
Unit=Mybackuper.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
